# shell/supervisor/elasticsearch.conf
[program:elasticsearch]
command=elasticsearch
directory=/opt/elasticsearch
autorestart=true
user=elasticsearch
