# shell/supervisor/kibana.conf
[program:kibana]
command=kibana
directory=/opt/kibana
autorestart=true
user=elasticsearch
