# shell/supervisor/logstash.conf
[program:logstash]
command=logstash
directory=/opt/logstash
autorestart=true
user=elasticsearch
